# infra/docker/vlc/config/onscreens-server.conf
[program:onscreens-server]
directory=/opt/tripbot
command=script/start-onscreens-server.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
startsecs=2
