# infra/docker/vlc/config/syslog.conf
[program:syslog]
command=/usr/sbin/syslog-ng -F
priority=1
autostart=true
autorestart=true
stdout_logfile=/var/log/syslog
stderr_logfile=/var/log/syslog
