# infra/docker/vlc/config/vlc.conf
[program:vlc]
directory=/opt/tripbot
command=script/x11/start-vlc.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
startsecs=2
