# infra/docker/vlc/config/vnc.conf
[program:vnc]
directory=/opt/tripbot
command=script/x11/start-vnc.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
