# infra/docker/obs/supervisor/admin-shim.conf
[program:admin-shim]
command=/opt/obs/script/start-admin-shim.sh
priority=10
autostart=true
# The shim is a single-page Flask app — Restarting it from a crash is
# always the right move (no external-state coupling). startretries is
# generous because losing the shim doesn't take the pod down (only
# OBS exiting does, per obs.conf), so transient flakes shouldn't fail
# the pod.
autorestart=true
startsecs=3
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
